# infra/services/besu-node-2.service
[Unit]
Description=Besu Node 2 - QBFT Validator
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/besu/Node-2
Environment="JAVA_HOME=/usr/lib/jvm/java-1.21.0-openjdk-amd64"
ExecStart=/home/ubuntu/besu/Node-2/run.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=besu-node-2

[Install]
WantedBy=multi-user.target
